# roles/xray/templates/xray-sync-webhook.service.j2
[Unit]
Description=Trigger Telegram bot Xray user sync when xray container starts
After=docker.service
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/xray-sync-webhook.sh
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
